# etc/supervisor/conf.d/supervisord.conf
# READONLY FILE, DO NOT EDIT THIS FILE
# conf/supervisord_expo.conf
[program:expo]
command=yarn expo start --tunnel --port 3000
directory=/app/frontend
autostart=true
autorestart=true
environment=EXPO_DEVTOOLS_LISTEN_ADDRESS="0.0.0.0",REACT_NATIVE_PACKAGER_HOSTNAME="0.0.0.0",CI="true"
stderr_logfile=/var/log/supervisor/expo.err.log
stdout_logfile=/var/log/supervisor/expo.out.log
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true

[program:backend]
command=/root/.venv/bin/uvicorn server:app --host 0.0.0.0 --port 8001 --workers 1 --reload
directory=/app/backend
autostart=true
autorestart=true
environment=APP_URL="https://b0b2b1b6-e5cd-4885-b126-d4286e9f5b07.preview.emergentagent.com",INTEGRATION_PROXY_URL="https://integrations.emergentagent.com"
stderr_logfile=/var/log/supervisor/backend.err.log
stdout_logfile=/var/log/supervisor/backend.out.log
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true

[program:mongodb]
command=/usr/bin/mongod --bind_ip_all
autostart=true
autorestart=true
stderr_logfile=/var/log/mongodb.err.log
stdout_logfile=/var/log/mongodb.out.log
